# etc/supervisord.d/galaxy.ini
[program:galaxy]
command = /opt/galaxy/run.sh
redirect_stderr=true
#stdout_logfile = /var/log/supervisor/%(program_name)s.log
stdout_logfile = /var/log/supervisor/%(program_name)s.log
